# BaukDockerSwarm/ExampleSwarm/images/jenkins-ssh_slave/setup.service
[Unit]
Description=SSH user and folder setup service

[Service]
ExecStart=/setup.sh

[Install]
WantedBy=multi-user.target
